# etc/systemd-upsource.service.in
#  This is part of the the systemd configuration for upsource

[Unit]
Description=Download or update local files remote remote sources based on sourcetabs in /etc/upsource.d/
Conflicts=shutdown.target
ConditionPathExists=/etc/upsource.d/config
ConditionPathExistsGlob=/etc/upsource.d/*.srctab
Requires=network.target local-fs.target
After=cloud-final.service
DefaultDependencies=false

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/upsource /etc/upsource.d/
#Environment=SYSTEMD_LOG_LEVEL=debug

[Install]
WantedBy=default.target
